# vdirsyncer/.config/systemd/user/vdirsyncer.service
[Unit]
Description=Calendar and contacts synchronization service

[Service]
Type=oneshot
ExecStart=-/usr/bin/vdirsyncer sync
