# redshift/.config/systemd/user/redshift.service.d/display.conf
[Service]
Environment=DISPLAY=:0
RestartSec=5s
Restart=on-failure
